# assets/pwsp-daemon.service
[Unit]
Description=Pipewire Soundpad Daemon
After=pipewire.service

[Service]
ExecStart=pwsp-daemon
Restart=no
RuntimeDirectory=pwsp

[Install]
WantedBy=default.target